# EliteSx2Mqtt/DEBIAN/postinst
#!/bin/bash

systemctl enable /opt/bin/elitesx2mqtt/elitesx2mqtt.service
systemctl start elitesx2mqtt.service
